# debian/postinst
#!/bin/sh

systemctl daemon-reload
systemctl enable driver-installer.service
systemctl start driver-installer.service

exit 0
